# _common/usr/lib/systemd/system/var-mnt.automount
[Unit]
Description=Automount orpheus:

[Automount]
Where=/var/mnt

[Install]
WantedBy=multi-user.target